# discovery/roles/monitor_thread/templates/duplicate_node_cleanup.service.j2
[Unit]
Description=Service to delete duplicate node

[Service]
Type=simple
ExecStart={{ duplicate_node_cleanup_shell_dest }}
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target