# roles/beegfs_role/files/beegfs-client.service.d/override.conf
# The reason for this change is to enable beegfs-client.service to persist
# through reboots in hyperconverged setup where the client and the oss server is
# running node and the client is not allowed to mount until oss server is ready.
[Service]
Type=forking
Restart=on-failure
RestartSec=5s
